# machine-4/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:mysql]
command=/usr/bin/mysqld_safe --datadir=/var/lib/mysql

[program:ssh]
command=/usr/sbin/sshd -D

[program:nginx]
command=nginx -g "daemon off;"
